# compiler/gravity/hemera-orbit-template.target.in
[Unit]
Description=Orbit @ORBIT@@DEBUG@
AllowIsolate=yes
