# compiler/gravity/hemera-star-template.target.in
[Unit]
Description=Hemera Star @NAME@
AllowIsolate=yes
